# scripts/systemd/belka-internal-service-ip.service
[Unit]
Description=BELKASQL internal floating service IP reconciliation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=CONFIG_FILE=/etc/belkasql/internal-service-ip.env
ExecStart=/usr/local/bin/belka-internal-service-ip.sh
